# kapowarr/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Nginx service run script for Kapowarr
# ==============================================================================

bashio::log.info "Starting nginx service..."

# Wait for Kapowarr backend to be ready on port 5656
bashio::log.info "Waiting for Kapowarr backend on port 5656..."
if ! bashio::net.wait_for 5656 localhost 300; then
    bashio::log.error "Kapowarr did not become ready on port 5656 within 300s"
    exit 1
fi

bashio::log.info "Kapowarr backend is ready"

if ! nginx -t 2>&1; then
    bashio::log.error "Nginx configuration test failed"
    exit 1
fi

bashio::log.info "Nginx starting in foreground mode"
exec nginx

