# cleanuparr/rootfs/etc/s6-overlay/s6-rc.d/nginx/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

declare -r exit_code="${1}"

if [[ "${exit_code}" -ne 0 ]] && [[ "${exit_code}" -ne 256 ]]; then
    bashio::log.error "nginx crashed (exit code ${exit_code}) — stopping add-on"
    bashio::log.error "Check the nginx error log output above for details"

    echo "${exit_code}" > /run/s6-linux-init-container-results/exitcode

    exec /run/s6/basedir/bin/halt
fi
