# .claude/skills/ha-apps/references/scaffold/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# services.d/nginx/run — nginx service run script
#
# s6-overlay starts nginx AFTER all cont-init.d scripts have completed.
# This script waits for the backend application to be ready before starting
# nginx, so the proxy does not return 502 errors during startup.
#
# The backend port is read from $APP_PORT, which is set in the Dockerfile
# ENV instruction and pre-populated into the container environment by s6-overlay.
# CUSTOMIZE: Adjust the wait_timeout below for apps with long startup times.
# ==============================================================================

# Redirect stderr to stdout so nginx errors appear in the HA add-on log
exec 2>&1

# Source the shared logging library
# shellcheck source=/usr/local/lib/ha-log.sh
source /usr/local/lib/ha-log.sh

ha::log::init "nginx"

log_info "Starting nginx ingress proxy..."

# ---------------------------------------------------------------------------
# Wait for the backend application to be ready
#
# bashio::net.wait_for <port> <host> <timeout_seconds>
# Returns 0 when a TCP connection to <host>:<port> succeeds within the timeout.
#
# APP_PORT is set in the Dockerfile ENV instruction and pre-populated by
# s6-overlay into the container environment.
# CUSTOMIZE: Adjust the timeout (default: 300s) for apps with long startup times.
# ---------------------------------------------------------------------------
readonly backend_host="localhost"
readonly wait_timeout=300

log_info "Waiting for APP_NAME backend on ${backend_host}:${APP_PORT} (timeout: ${wait_timeout}s)..."

if ! bashio::net.wait_for "${APP_PORT}" "${backend_host}" "${wait_timeout}"; then
    log_error "APP_NAME backend did not become ready within ${wait_timeout}s"
    log_error "Check the APP_NAME service logs for startup errors"
    exit 1
fi

log_info "APP_NAME backend is ready on ${backend_host}:${APP_PORT}"

# ---------------------------------------------------------------------------
# Final nginx configuration check before starting
# ---------------------------------------------------------------------------
if ! nginx -t 2>&1; then
    log_error "nginx configuration test failed — check /etc/nginx/servers/ingress.conf"
    exit 1
fi

# ---------------------------------------------------------------------------
# Start nginx in the foreground.
# "daemon off;" in nginx.conf ensures nginx does not fork into the background.
# 'exec' replaces this shell with the nginx process so s6 can supervise it.
# ---------------------------------------------------------------------------
log_info "nginx starting in foreground mode"
exec nginx
