# kapowarr/rootfs/etc/services.d/kapowarr/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Handle Kapowarr service exit/crash
# ==============================================================================

# Exit code 0 = normal shutdown, 256 = killed by signal, other = crash
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
    bashio::log.error "Kapowarr crashed with exit code ${1}, stopping add-on"
    s6-svscanctl -t /var/run/s6/services
fi

