# kapowarr/rootfs/etc/services.d/kapowarr/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# services.d/kapowarr/run — Kapowarr service run script
#
# s6-overlay v3 calls this script to START the Kapowarr application service.
# The script MUST run the application in the FOREGROUND.
# ==============================================================================

# Redirect stderr to stdout
exec 2>&1

# Source the shared logging library
# shellcheck source=/usr/local/lib/ha-log.sh
source /usr/local/lib/ha-log.sh

ha::log::init "kapowarr"

log_info "Starting Kapowarr..."

# ---------------------------------------------------------------------------
# Read configuration
# ---------------------------------------------------------------------------
declare log_level
log_level="$(bashio::config 'log_level' 'info')"
log_info "Log level: ${log_level}"

# ---------------------------------------------------------------------------
# Set runtime environment variables
# ---------------------------------------------------------------------------
export TZ="${TZ:-UTC}"

# Get the port from environment or default to 5656
declare app_port
app_port="${APP_PORT:-5656}"

log_info "Kapowarr will listen on 0.0.0.0:${app_port}"

# ---------------------------------------------------------------------------
# Change to the application directory
# ---------------------------------------------------------------------------
cd /app || bashio::exit.nok "Could not change to /app directory"

# ---------------------------------------------------------------------------
# Start Kapowarr
# Kapowarr is a Python Flask application. The entry point is Kapowarr.py
# It reads configuration from environment variables and command-line args.
# ---------------------------------------------------------------------------
exec python3 /app/Kapowarr.py
