# kapowarr/rootfs/etc/s6-overlay/s6-rc.d/kapowarr/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# s6-rc.d/kapowarr/finish — Application service finish script
#
# s6-overlay calls this script when the 'run' script exits for any reason.
#
# Arguments passed by s6:
#   $1 — exit code of the run script
#          0   = clean shutdown (app exited normally)
#          256 = app was killed by a signal
#          N   = app exited with error code N
#
# Kapowarr restart codes (since we bypass its subprocess wrapper):
#   131 = RestartVersion.NORMAL — app requested a normal restart
#   132 = RestartVersion.HOSTING_CHANGES — hosting settings changed
#
# Behaviour:
#   - Normal exit (0) or signal kill (256): s6 restarts automatically
#   - Kapowarr restart codes (131, 132): s6 restarts automatically
#   - Any other exit code: treat as a crash, stop the add-on
# ==============================================================================

declare -r exit_code="${1}"

bashio::log.info "Kapowarr exited with code ${exit_code}"

# Allow clean exits, signal kills, and Kapowarr restart requests
if [[ "${exit_code}" -eq 0 ]] \
    || [[ "${exit_code}" -eq 256 ]] \
    || [[ "${exit_code}" -eq 131 ]] \
    || [[ "${exit_code}" -eq 132 ]]; then
    bashio::log.info "Kapowarr will restart automatically"
    exit 0
fi

# Any other exit code is a crash — stop the add-on
bashio::log.error "Kapowarr crashed (exit code ${exit_code}) — stopping add-on to prevent restart loop"
bashio::log.error "Check the logs above for the cause of the crash"

echo "${exit_code}" > /run/s6-linux-init-container-results/exitcode
exec /run/s6/basedir/bin/halt
