# kapowarr/rootfs/etc/s6-overlay/s6-rc.d/kapowarr/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# s6-rc.d/kapowarr/run — Kapowarr service run script
#
# s6-overlay v3 calls this script to START the Kapowarr application service.
# The script MUST run the application in the FOREGROUND.
# ==============================================================================

# Redirect stderr to stdout
exec 2>&1

bashio::log.info "Starting Kapowarr..."

# ---------------------------------------------------------------------------
# Read configuration from s6 container environment (exported by kapowarr-setup)
# ---------------------------------------------------------------------------
declare log_level="${APP_LOG_LEVEL:-info}"
bashio::log.level "${log_level}"
bashio::log.info "Log level: ${log_level}"

# ---------------------------------------------------------------------------
# Set runtime environment variables
# ---------------------------------------------------------------------------
export TZ="${TZ:-UTC}"

# Get the port from environment or default to 5656
declare app_port
app_port="${APP_PORT:-5656}"

bashio::log.info "Kapowarr will listen on 0.0.0.0:${app_port}"

# ---------------------------------------------------------------------------
# Resolve HA ingress path for Kapowarr's native url_base setting
# ---------------------------------------------------------------------------
# Kapowarr supports a url_base setting that wraps the Flask app with
# DispatcherMiddleware, making all routes, static assets, API endpoints,
# and WebSocket paths respect the prefix natively. This eliminates the
# need for nginx sub_filter rewriting.
#
# The ingress path (e.g. /api/hassio_ingress/<token>) is stable across
# container restarts and only changes on add-on reinstall.
declare ingress_path
ingress_path="$(bashio::addon.ingress_entry)"
bashio::log.info "Ingress path (url_base): ${ingress_path}"

# ---------------------------------------------------------------------------
# Change to the application directory
# ---------------------------------------------------------------------------
cd /app || bashio::exit.nok "Could not change to /app directory"

# ---------------------------------------------------------------------------
# Map HA log level string to Python logging integer
# ---------------------------------------------------------------------------
declare log_level_int
case "${log_level}" in
    trace)   log_level_int=5 ;;
    debug)   log_level_int=10 ;;
    info)    log_level_int=20 ;;
    warning) log_level_int=30 ;;
    error)   log_level_int=40 ;;
    *)       log_level_int=20 ;;
esac

# ---------------------------------------------------------------------------
# Pre-seed database with url_base BEFORE starting Kapowarr
# ---------------------------------------------------------------------------
# Kapowarr reads url_base from the config table at startup and caches it
# in a Python Singleton. Direct sqlite3 changes after startup are invisible
# to the running process. Therefore, the value MUST be in the DB before
# Kapowarr starts.
#
# This works because:
#   - Kapowarr's setup_db() uses CREATE TABLE IF NOT EXISTS — our
#     pre-created table is preserved
#   - Settings._insert_missing_settings() uses INSERT OR IGNORE — our
#     pre-seeded values are never overwritten by defaults
#   - Settings._fetch_settings() reads our correct url_base
#   - SERVER.set_url_base() configures DispatcherMiddleware correctly
# ---------------------------------------------------------------------------
declare db_file="/config/data/Kapowarr.db"

sqlite3 "${db_file}" <<SQL
CREATE TABLE IF NOT EXISTS config(
    key VARCHAR(100) PRIMARY KEY,
    value BLOB
);
INSERT OR REPLACE INTO config(key, value) VALUES ('url_base', '${ingress_path}');
INSERT OR REPLACE INTO config(key, value) VALUES ('log_level', ${log_level_int});
SQL

bashio::log.info "Database pre-seeded: url_base=${ingress_path}, log_level=${log_level_int}"

# ---------------------------------------------------------------------------
# Update PWA manifest with the correct url_base prefix
# ---------------------------------------------------------------------------
# Kapowarr's update_manifest() is normally called only when url_base changes
# via the Settings API. Since we set url_base directly in the DB, we must
# update the manifest file ourselves.
python3 -c "
import json, sys
url_base = sys.argv[1]
path = 'frontend/static/json/pwa_manifest.json'
with open(path, 'r') as f:
    m = json.load(f)
m['start_url'] = url_base + '/'
m['icons'][0]['src'] = url_base + '/static/img/favicon.svg'
with open(path, 'w') as f:
    json.dump(m, f, indent=4)
" "${ingress_path}" 2>/dev/null || bashio::log.warning "Could not update PWA manifest"

# ---------------------------------------------------------------------------
# Start Kapowarr in the foreground with exec for proper signal handling
# ---------------------------------------------------------------------------
exec python3 /app/Kapowarr.py -d /config/data
