# kapowarr/rootfs/etc/s6-overlay/s6-rc.d/kapowarr/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# s6-rc.d/kapowarr/run — Kapowarr service run script
#
# s6-overlay v3 calls this script to START the Kapowarr application service.
# The script MUST run the application in the FOREGROUND.
# ==============================================================================

# Redirect stderr to stdout
exec 2>&1

bashio::log.info "Starting Kapowarr..."

# ---------------------------------------------------------------------------
# Read configuration from s6 container environment (exported by kapowarr-setup)
# ---------------------------------------------------------------------------
declare log_level="${APP_LOG_LEVEL:-info}"
bashio::log.level "${log_level}"

export TZ="${TZ:-UTC}"

declare app_port="${APP_PORT:-5656}"

# Read cached ingress path from s6 environment (set by kapowarr-setup)
declare ingress_path="${APP_INGRESS_ENTRY}"
bashio::log.info "Ingress url_base: ${ingress_path}"

# ---------------------------------------------------------------------------
# Change to the application directory
# ---------------------------------------------------------------------------
cd /app || bashio::exit.nok "Could not change to /app directory"

# ---------------------------------------------------------------------------
# Map HA log level string to Python logging integer
# ---------------------------------------------------------------------------
declare log_level_int
case "${log_level}" in
    trace)   log_level_int=5 ;;
    debug)   log_level_int=10 ;;
    info)    log_level_int=20 ;;
    warning) log_level_int=30 ;;
    error)   log_level_int=40 ;;
    *)       log_level_int=20 ;;
esac

# ---------------------------------------------------------------------------
# Pre-seed database with url_base BEFORE starting Kapowarr
# ---------------------------------------------------------------------------
# Kapowarr reads url_base from the config table at startup and caches it
# in a Python Singleton. Direct sqlite3 changes after startup are invisible
# to the running process. Therefore, the value MUST be in the DB before
# Kapowarr starts.
#
# This works because:
#   - Kapowarr's setup_db() uses CREATE TABLE IF NOT EXISTS — our
#     pre-created table is preserved
#   - Settings._insert_missing_settings() uses INSERT OR IGNORE — our
#     pre-seeded values are never overwritten by defaults
#   - Settings._fetch_settings() reads our correct url_base
#   - SERVER.set_url_base() configures DispatcherMiddleware correctly
# ---------------------------------------------------------------------------
declare db_file="/config/data/Kapowarr.db"

sqlite3 "${db_file}" <<SQL
CREATE TABLE IF NOT EXISTS config(
    key VARCHAR(100) PRIMARY KEY,
    value BLOB
);
INSERT OR REPLACE INTO config(key, value) VALUES ('url_base', '${ingress_path}');
INSERT OR REPLACE INTO config(key, value) VALUES ('log_level', ${log_level_int});
SQL

bashio::log.debug "Database pre-seeded: url_base=${ingress_path}, log_level=${log_level_int}"

# ---------------------------------------------------------------------------
# Update PWA manifest with the correct url_base prefix
# ---------------------------------------------------------------------------
# Kapowarr's update_manifest() is normally called only when url_base changes
# via the Settings API. Since we set url_base directly in the DB, we must
# update the manifest file ourselves. Uses sed instead of Python to avoid
# starting a full interpreter for a trivial JSON edit.
declare manifest="frontend/static/json/pwa_manifest.json"
if [[ -f "${manifest}" ]]; then
    sed -i \
        -e "s|\"start_url\": \"[^\"]*\"|\"start_url\": \"${ingress_path}/\"|" \
        -e "s|\"src\": \"[^\"]*favicon.svg\"|\"src\": \"${ingress_path}/static/img/favicon.svg\"|" \
        "${manifest}"
fi

# ---------------------------------------------------------------------------
# Start Kapowarr — bypass the subprocess wrapper for faster startup
# ---------------------------------------------------------------------------
# Kapowarr normally uses a parent→child subprocess pattern (for internal
# restart support). The parent process imports ALL modules, then spawns a
# child via Popen that imports them ALL AGAIN — doubling Python startup time.
#
# By setting KAPOWARR_RUN_MAIN=1, we skip the parent process entirely and
# go straight to _main(). This eliminates one full Python startup cycle
# (~3-6s savings depending on hardware).
#
# Internal restarts (exit codes 131/132) are handled by the finish script,
# which lets s6 restart this run script instead.
#
# Signal handling: Python's default SIGTERM handler terminates without
# running cleanup. We run Python in the background and trap SIGTERM,
# forwarding it as SIGINT (which raises KeyboardInterrupt and triggers
# the try/finally cleanup in _main()).
# ---------------------------------------------------------------------------
export KAPOWARR_RUN_MAIN=1

python3 -u /app/Kapowarr.py -d /config/data &
pid=$!

# Forward SIGTERM as SIGINT for graceful Python shutdown
trap 'kill -INT ${pid} 2>/dev/null' TERM

wait "${pid}"
exit $?
