# .claude/skills/ha-apps3/references/scaffold/rootfs/etc/services.d/APP_NAME/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# services.d/APP_NAME/run — Application service run script
#
# s6-overlay v3 calls this script to START the application service.
# The script MUST run the application in the FOREGROUND (no daemon mode, no &).
# s6 supervises this process and automatically restarts it if it exits.
#
# "with-contenv" makes variables written to /var/run/s6/container_environment/
# by the cont-init.d scripts available here as normal shell variables.
#
# CUSTOMIZE: Replace the exec line at the bottom with your application's
#            startup command.
# ==============================================================================

# Redirect stderr to stdout so all output appears in the HA add-on log panel
exec 2>&1

# Source the shared logging library
# shellcheck source=/usr/local/lib/ha-log.sh
source /usr/local/lib/ha-log.sh

ha::log::init "APP_NAME"

log_info "Starting APP_NAME..."

# ---------------------------------------------------------------------------
# Read configuration
# CUSTOMIZE: read any additional options you need here
# ---------------------------------------------------------------------------
declare log_level
log_level="$(bashio::config 'log_level' 'info')"
log_info "Log level: ${log_level}"

# ---------------------------------------------------------------------------
# Set runtime environment variables
# CUSTOMIZE: export any environment variables your application needs
# ---------------------------------------------------------------------------
export TZ="${TZ:-UTC}"
export APP_LOG_LEVEL="${log_level}"

# CUSTOMIZE: most applications listen on a fixed internal port.
# This should match the port in servers/ingress.conf proxy_pass.
readonly APP_PORT=APP_PORT
log_info "APP_NAME will listen on 0.0.0.0:${APP_PORT}"

# ---------------------------------------------------------------------------
# Change to the application directory
# CUSTOMIZE: update the path to match where the app is installed
# ---------------------------------------------------------------------------
cd /app || bashio::exit.nok "Could not change to /app directory"

# ---------------------------------------------------------------------------
# IMPORTANT: Use 'exec' so the application process REPLACES this shell.
# This ensures s6 correctly tracks the PID and signal delivery works.
#
# CUSTOMIZE: Replace the command below with your application's start command.
# Examples:
#   exec ./myapp --config /config/myapp.yaml
#   exec python3 main.py
#   exec node server.js
#   exec java -jar app.jar
# ---------------------------------------------------------------------------
exec /usr/bin/APP_NAME \
    --host 0.0.0.0 \
    --port "${APP_PORT}" \
    --config /config
