# .claude/skills/ha-apps/references/scaffold/rootfs/etc/s6-overlay/s6-rc.d/APP_NAME/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# s6-rc.d/APP_NAME/finish — Application service finish script
#
# s6-overlay calls this script when the 'run' script exits for any reason.
#
# Arguments passed by s6:
#   $1 — exit code of the run script
#          0   = clean shutdown (app exited normally)
#          256 = app was killed by a signal
#          N   = app exited with error code N
#
# Behaviour implemented here:
#   - Normal exit (0) or signal kill (256): do nothing, s6 restarts automatically
#   - Any other exit code: treat as a crash and stop the entire add-on
#
# The crash-stop behaviour prevents restart loops that would fill the log and
# consume resources.  The user must check the logs to diagnose the issue.
# ==============================================================================

declare -r exit_code="${1}"

# Source the full HA framework (loads all libraries including ha-log.sh)
# shellcheck source=/usr/local/lib/ha-framework/ha-framework.sh
source /usr/local/lib/ha-framework/ha-framework.sh

ha::log::init "APP_NAME"

log_info "APP_NAME exited with code ${exit_code}"

if [[ "${exit_code}" -ne 0 ]] && [[ "${exit_code}" -ne 256 ]]; then
    log_error "APP_NAME crashed (exit code ${exit_code}) — stopping add-on to prevent restart loop"
    log_error "Check the logs above for the cause of the crash"

    # Write the exit code so the container returns a non-zero status
    echo "${exit_code}" > /run/s6-linux-init-container-results/exitcode

    # Halt the entire s6 supervision tree (stops the add-on)
    exec /run/s6/basedir/bin/halt
fi

# Codes 0 and 256 fall through here.
# s6 will automatically restart the 'run' script after the configured backoff.
log_info "APP_NAME will restart automatically"
