# cleanuparr/rootfs/etc/s6-overlay/s6-rc.d/cleanuparr-setup/up
/etc/s6-overlay/scripts/cleanuparr-setup
