# kapowarr/rootfs/etc/s6-overlay/s6-rc.d/nginx/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# s6-rc.d/nginx/finish — nginx service finish script
#
# s6-overlay calls this script when the nginx 'run' script exits.
#
# Arguments passed by s6:
#   $1 — exit code of the run script
#          0   = clean shutdown
#          256 = nginx was killed by a signal
#          N   = nginx exited with error code N
#
# If nginx exits abnormally it indicates a configuration error or resource
# problem that requires operator intervention.  We halt the add-on rather
# than looping endlessly.
# ==============================================================================

declare -r exit_code="${1}"

# Source the HA framework (loads all libraries)
# shellcheck source=/usr/local/lib/ha-framework/ha-framework.sh
source /usr/local/lib/ha-framework/ha-framework.sh

ha::log::init "nginx"

if [[ "${exit_code}" -ne 0 ]] && [[ "${exit_code}" -ne 256 ]]; then
    log_error "nginx crashed (exit code ${exit_code}) — stopping add-on"
    log_error "Check the nginx error log output above for details"

    # Write the exit code so the container returns a non-zero status
    echo "${exit_code}" > /run/s6-linux-init-container-results/exitcode

    # Halt the entire s6 supervision tree (stops the add-on)
    exec /run/s6/basedir/bin/halt
fi
