# cleanuparr/rootfs/etc/s6-overlay/s6-rc.d/cleanuparr/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

declare -r exit_code="${1}"

bashio::log.info "Cleanuparr exited with code ${exit_code}"

if [[ "${exit_code}" -ne 0 ]] && [[ "${exit_code}" -ne 256 ]]; then
    bashio::log.error "Cleanuparr crashed (exit code ${exit_code}) — stopping add-on to prevent restart loop"
    bashio::log.error "Check the logs above for the cause of the crash"

    echo "${exit_code}" > /run/s6-linux-init-container-results/exitcode

    exec /run/s6/basedir/bin/halt
fi

bashio::log.info "Cleanuparr will restart automatically"
