# cleanuparr/rootfs/etc/s6-overlay/s6-rc.d/cleanuparr/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

set -e

exec 2>&1

bashio::log.info "Starting Cleanuparr..."

LOG_LEVEL=$(bashio::config 'log_level')
DRY_RUN=$(bashio::config 'dry_run')

bashio::log.info "Log level: ${LOG_LEVEL}"
bashio::log.info "Dry run: ${DRY_RUN}"

mkdir -p /config

export ASPNETCORE_URLS="http://0.0.0.0:11011"
export ASPNETCORE_ENVIRONMENT="Production"
export LOGLEVEL__DEFAULT="${LOG_LEVEL}"

cd /app || bashio::exit.nok "Could not change to /app directory"

exec /usr/bin/cleanuparr
